# debian/fusioninventory-agent.cron.weekly
#!/bin/sh

if [ -f /etc/default/fusioninventory-agent ] ; then
	. /etc/default/fusioninventory-agent
fi

if [ ! "$UPDATE_IDS" = "yes"  ]; then
	exit 0
fi

PCIIDS_UPDATE=/usr/sbin/update-pciids
USBIDS_UPDATE=/usr/sbin/update-usbids

for PROG in $PCIIDS_UPDATE $USBIDS_UPDATE
do
	if [ ! -x "$PROG" ]; then
		continue
	fi

	$PROG >/dev/null 2>&1
done
